# urmart-demo/crontab
0 * * * * /usr/bin/python3 /workspace/send_shop_report.py >> /sent.log 2>&1
